# roles/tomcat_server/templates/tomcat.service.j2
# Systemd unit file for tomcat
[Unit]
Description=Apache Tomcat Web Application Container
After=syslog.target network.target

[Service]
Type=forking

Environment="JRE_HOME=/usr/lib/jvm/jre"
Environment="CLASSPATH=.:$JRE_HOME/lib/ext"
Environment="INSTANCE_OWNER={{ tomcat_user.name }}"
Environment="CATALINA_HOME=/opt/tomcat"
Environment="CATALINA_BASE=/home/$INSTANCE_OWNER/tomcat"
Environment="CATALINA_PID=/home/$INSTANCE_OWNER/tomcat/run/${INSTANCE_OWNER}.pid"

ExecStart=/home/{{ tomcat_user.name }}/tomcat/bin/startup.sh
ExecStop=/home/{{ tomcat_user.name }}/tomcat/bin/shutdown.sh

User={{ tomcat_user.name }}
Group={{ tomcat_group.name }}
Restart=on-failure

[Install]
WantedBy=multi-user.target