# roles/tomcat_server/templates/tomcat.service.j2
# Systemd unit file for tomcat
[Unit]
Description=Apache Tomcat Web Application Container
After=syslog.target network.target

[Service]
Type=forking

WorkingDirectory=/home/{{ tomcat_user.name }}

Environment="INSTANCE_OWNER={{ tomcat_user.name }}"
Environment="CATALINA_HOME=/opt/tomcat"
Environment="CATALINA_BASE=/home/$INSTANCE_OWNER/tomcat"

ExecStart=/home/{{ tomcat_user.name }}/tomcat/bin/startup.sh
ExecStop=/home/{{ tomcat_user.name }}/tomcat/bin/shutdown.sh

User={{ tomcat_user.name }}
Group={{ tomcat_group.name }}
RestartSec=10
Restart=on-failure

[Install]
WantedBy=multi-user.target